# warema-bridge/etc/services.d/warema-bridge/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Home Assistant Community Add-on: Warema Bridge
# Take down the S6 supervision tree when the Warema bridge fails
# ==============================================================================
bashio::log.info "Stopping the Warema bridge..."

# NodeJS wird automatisch SIGTERM erhalten, s6 sorgt für sauberen Stop
